# change-csv.service
[Unit]
Description= Copier les fichiers CSV
[Service]
Type=oneshot
ExecStart=/path/to/the/same/directory/as/main.py/change-csv.sh #change it

